# supervisor/supervisord.d/stockvip_celery_worker.conf
[program:stockvip-celery-worker]
command=/root/.virtualenvs/StockVip/bin/python  /root/code/StockVip/manage.py celery worker -l INFO
directory=/root/code/StockVip/
environment=C_FORCE_ROOT="true"
umask=022
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/root/code/StockVip/logs/stockvip-celery-worker.log
stderr_logfile=/root/code/StockVip/logs/stockvip-celery-workererr.log
stopsignal=QUIT
killasgroup=true
